# docker/rootfs/etc/s6-overlay/s6-rc.d/generate-env/up
#!/usr/bin/with-contenv sh
exec /etc/s6-overlay/scripts/generate-env.sh
